# crontab
* * * * * /root/scripts/nfs_disk_usage.sh $DIRECTORY
